# package/docker/supervisor_conf.d/mongod.conf
[program:mongod]
command=/usr/bin/mongod

