# package/docker/supervisor_conf.d/nginx.conf
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"


