# package/docker/supervisor_conf.d/redis-server.conf
[program:redis-server]
command=/usr/bin/redis-server

