# package/docker/supervisor_conf.d/sharelatex-docstore.conf
[program:sharelatex-docstore]
command=/usr/bin/node /var/www/sharelatex/docstore/app.js

