# package/docker/supervisor_conf.d/sharelatex-chat.conf
[program:sharelatex-chat]
command=/usr/bin/node /var/www/sharelatex/chat/app.js

