# package/docker/supervisor_conf.d/sharelatex-filestore.conf
[program:sharelatex-clsi]
command=/usr/bin/node /var/www/sharelatex/clsi/app.js

