# package/docker/supervisor_conf.d/sharelatex-spelling.conf
[program:sharelatex-spelling]
command=/usr/bin/node /var/www/sharelatex/tags/spelling.js

