# package/docker/supervisor_conf.d/sharelatex-web.conf
[program:sharelatex-web]
command=/usr/bin/node /var/www/sharelatex/web/app.js

