# package/docker/supervisor_conf.d/sharelatex-tags.conf
[program:sharelatex-tags]
command=/usr/bin/node /var/www/sharelatex/tags/app.js

