# package/docker/supervisor_conf.d/sharelatex-track-changes.conf
[program:sharelatex-track-changes]
command=/usr/bin/node /var/www/sharelatex/track-changes/app.js

