# package/docker/supervisor_conf.d/sharelatex-document-updater.conf
[program:sharelatex-document-updater]
command=/usr/bin/node /var/www/sharelatex/document-updater/app.js

