# server-ce/runit/real-time-overleaf/run
#!/bin/bash
exec /sbin/setuser www-data /usr/bin/node /overleaf/services/real-time/app.js >> /var/log/overleaf/real-time.log 2>&1
