# runit/errorlogs/run
#!/bin/bash
tail -F /var/log/sharelatex/*.log | grep '"level":50' >> /var/log/sharlatex-errors.log
